# battery-notify.service
[Service]
ExecStart=battery-notify
Restart=always

[Install]
WantedBy=default.target
